# deploy/fluxrouter.service
[Unit]
Description=FluxRouter - AI API Gateway Platform
Documentation=https://github.com/Wei-Shaw/fluxrouter
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=fluxrouter
Group=fluxrouter
WorkingDirectory=/opt/fluxrouter
ExecStart=/opt/fluxrouter/fluxrouter
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fluxrouter

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/fluxrouter

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
